# systemd/user/xdg-desktop-portal.service.d/override.conf
[Unit]
Requisite=
PartOf=
After=dbus.socket basic.target session.slice
